# deploy/systemd/agent-reactor-runtime.service
# agent-reactor runtime — pty session host. Long-lived; owns the session
# socket and the on-disk session store. Web-server (`agent-reactor-server`)
# and any future native client attach to this unit's socket; restarting them
# does not lose sessions. Restarting THIS unit drops every attached client
# and recreates sessions from disk on the next boot.
#
# Run as a per-user unit: `systemctl --user enable --now agent-reactor-web`
# cascades runtime + server via Requires=. For boot-time autostart even when
# the user is not logged in: `loginctl enable-linger $USER`.
#
# Prerequisite: the invoking user must reach the Docker daemon — either by
# membership in the `docker` group OR by a rootless docker socket at
# $XDG_RUNTIME_DIR/docker.sock (auto-detected). Without that the runtime can
# still serve unsandboxed sessions but devcontainer launches will fail.

[Unit]
Description=agent-reactor runtime (session host)
# Operator guide: docs/user/systemd.md in the agent-reactor source repo.

[Service]
Type=simple
# StateDirectory creates %S/agent-reactor (XDG_STATE_HOME/agent-reactor,
# typically ~/.local/state/agent-reactor) with mode 0700 on first start. The
# runtime writes its socket, pid file, sessions snapshot, event log and
# runtime log all under this single dir today; splitting socket onto tmpfs
# (%t) is a future polish that needs a runtime -socket-path flag.
StateDirectory=agent-reactor
ExecStart=%h/.local/bin/agent-reactor-runtime -data-dir %S/agent-reactor
# reactor-prod keeps the docker container/label namespace distinct from the
# TUI daemon's `reactor` and run-dev.sh's `reactor-dev` — without this, mount
# hash drift would `docker rm -f` a peer daemon's container.
Environment=ROOST_DEVCONTAINER_PREFIX=reactor-prod
# Optional operator override file. Leading `-` makes systemd tolerate it
# being absent on first boot. Put TLS path overrides, alternate prefix etc.
EnvironmentFile=-%S/agent-reactor/runtime.env

Restart=on-failure
RestartSec=2s
KillSignal=SIGTERM
# Coordinator's graceful shutdown (signal handler in coordinator.go) reaps
# session goroutines and flushes the sessions snapshot. Observed <5s on a
# warm cache; 15s leaves headroom for a slow disk under heavy session count.
TimeoutStopSec=15s

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
